# central-api/com.jarvis.central-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.central-api</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/gunicorn</string>
        <string>--bind</string>
        <string>127.0.0.1:3003</string>
        <string>--workers</string>
        <string>2</string>
        <string>--timeout</string>
        <string>120</string>
        <string>--access-logfile</string>
        <string>/Users/clawdbot/clawd/central-api/logs/access.log</string>
        <string>--error-logfile</string>
        <string>/Users/clawdbot/clawd/central-api/logs/error.log</string>
        <string>--log-level</string>
        <string>info</string>
        <string>api.app:app</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/clawdbot/clawd/central-api</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/clawdbot/clawd/central-api</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/clawdbot/clawd/central-api/logs/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/clawdbot/clawd/central-api/logs/stderr.log</string>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>Nice</key>
    <integer>0</integer>
</dict>
</plist>
